# recipient_filter/recipient_filter.service
[Unit]
Description=Email policy server that validates recipient addresses with a secret
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/recipient_filter

[Install]
WantedBy=multi-user.target
